# deploy/supervisor/giftlab-reverb.conf
; Laravel Reverb websocket server - the ONLY real-time transport (no polling).
; Bound to localhost; Nginx (reverb.giftlab.example) terminates TLS and proxies
; the wss upgrade to it. Kept alive perpetually by Supervisor.
; Install to /etc/supervisor/conf.d/giftlab-reverb.conf
[program:giftlab-reverb]
process_name=%(program_name)s
command=php /var/www/giftlab/current/artisan reverb:start --host=127.0.0.1 --port=8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/giftlab/reverb.log
stopwaitsecs=10
